# systemd/myriad.service
[Unit]
Description=Boot Myriad devices
Before=kubelet.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/myriad
